# start-flask-app.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/obs/start-flask-app.sh

[Install]
WantedBy=multi-user.target
